# deploy/blog-2025.service
[Unit]
Description=2025 Blog Next.js server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/2025-blog-public/current
Environment=NODE_ENV=production
Environment=NEXT_TELEMETRY_DISABLED=1
Environment=NEXT_OUTPUT_MODE=server
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3001
ExecStart=/usr/bin/node /opt/2025-blog-public/current/node_modules/next/dist/bin/next start -H 0.0.0.0 -p 3001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
