# deploy/blog-2025.service
[Unit]
Description=2025 Blog Next.js standalone server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/2025-blog-public/current
Environment=NODE_ENV=production
Environment=NEXT_TELEMETRY_DISABLED=1
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3001
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
